# otto-lead-scraper.timer
[Unit]
Description=Otto Web Assist Lead Scraper — hourly
Requires=otto-lead-scraper.service

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
AccuracySec=5min

[Install]
WantedBy=timers.target
