# otto-reflection.timer
[Unit]
Description=Otto Reflection Timer (Half-Hour Mark)

[Timer]
OnCalendar=*:30
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
